# roles/pcmanfm/files/pcmanfm.service
[Unit]
Description=PCManFM daemon
PartOf=xsession.target

[Service]
ExecStart=/usr/bin/pcmanfm --daemon-mode

[Install]
WantedBy=xsession.target
