# roles/mpd/files/mpd.service.d/10-mpd-share-dependency.conf
[Unit]
After=
After=network.target sound.target mnt-alhena-music.mount
