# roles/mpd/files/mnt-share-music.automount
[Unit]
Description=Mount share music automount
Requires=network-online.target
After=network-online.service

[Automount]
Where=/mnt/share/music

[Install]
WantedBy=multi-user.target
