# roles/libinput-touchpad/files/libinput-gestures.service
[Unit]
Description=Translates trackpad gestures

[Service]
ExecStart=/usr/bin/libinput-gestures
Restart=always

[Install]
WantedBy=graphical.target
