# roles/X11/files/xsession.target
[Unit]
Description=Xsession is running
BindsTo=graphical-session.target
