# roles/redshift/files/geoclue-agent.service
[Unit]
Description=geoclue agent
PartOf=redshift.service

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
